# org.halext.api.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>org.halext.api</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/scawful/Code/halext-org/backend/env/bin/uvicorn</string>
        <string>main:app</string>
        <string>--host</string>
        <string>127.0.0.1</string>
        <string>--port</string>
        <string>8000</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/scawful/Code/halext-org/backend</string>

    <key>EnvironmentVariables</key>
    <dict>
        <!-- Removed DATABASE_URL to default to SQLite for development -->
        <!-- For production PostgreSQL, set: postgresql://halext_user:REDACTED_PASSWORD@localhost/halext_org -->
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/scawful/Code/halext-org/service.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/scawful/Code/halext-org/service.err.log</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>WatchPaths</key>
    <array>
        <string>/Users/scawful/Code/halext-org/backend/</string>
    </array>

</dict>
</plist>
