# infra/ubuntu/halext-api.service
[Unit]
Description=Halext FastAPI service
After=network.target

[Service]
WorkingDirectory=/srv/halext/halext-org-project/backend
EnvironmentFile=/srv/halext/halext-org-project/backend/.env
ExecStart=/srv/halext/halext-org-project/backend/env/bin/uvicorn main:app --host 127.0.0.1 --port 8000
Restart=always
User=www-data
Group=www-data

[Install]
WantedBy=multi-user.target
